# crontab
1 17-23 * * * killall -STOP qemu-system-x86_64

